# home/bin/work_log/logStartup.service
[Unit]
Description=Loggin Startups and shutdowns

[Service]
ExecStart=/home/mweber/bin/work_log/log.py Boot

[Install]
WantedBy=multi-user.target
